# files/systemd/sandbox.tractstack-tractstack.service
[Unit]
Description="Run tractstack build process"

[Service]
User=t8k
ExecStart=/home/t8k/sandbox/sandbox_ZZZZZ/scripts/build.sh
