# files/systemd/sandbox.tractstack-tractstack.path
[Unit]
Description="Trigger build process for Tract Stack"

[Path]
PathModified=/home/t8k/sandbox/ZZZZZ/releases/watch
Unit=t8k-sandbox_ZZZZZ.service

[Install]
WantedBy=multi-user.target
