# files/systemd/tractstack-tractstack.service
[Unit]
Description="Run tractstack build process"

[Service]
WorkingDirectory=/home/ZZZZZ/scripts
ExecStart=/home/ZZZZZ/scripts/build.sh
