# files/systemd/tractstack-tractstack.service
[Unit]
Description="Run tractstack build process"

[Service]
User=ZZZZZ
ExecStart=/home/ZZZZZ/scripts/build.sh
