# files/rsnapshot/systemd/t8k-backup-hourly.service
[Unit]
Description=t8k hourly backup service
[Service]
Type=oneshot
ExecStart=/usr/bin/rsnapshot hourly
