# files/rsnapshot/systemd/t8k-backup.service
[Unit]
Description=t8k backup service
[Service]
Type=oneshot
ExecStart=/usr/bin/rsnapshot daily
