# files/rsnapshot/systemd/t8k-mysql-backup.service
[Unit]
Description=t8k mysql backup service
[Service]
Type=oneshot
ExecStart=/home/t8k/tractstack-installer/scripts/t8k-mysql-backup.sh
