# files/rsnapshot/systemd/t8k-backup-hourly.timer
[Unit]
Description=t8k hourly backup timer
[Timer]
OnCalendar=hourly
Persistent=true
[Install]
WantedBy=timers.target
